# schedule/com.cuz.distill.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cuz.distill</string>

  <key>ProgramArguments</key>
  <array>
    <string>__BASH__</string>
    <string>__ENGINE__/bin/cuz-distill-prepare.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CUZ_CONFIG</key>
    <string>__CONFIG__</string>
    <key>CUZ_DISTILL_MODEL</key>
    <string>claude-sonnet-4-6</string>
    <key>CUZ_DISTILL_USD</key>
    <string>2.00</string>
  </dict>

  <!-- Weekly: Sunday 05:00. Prepares a proposal; never writes the library. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>5</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__CACHE__/launchd-distill.out.log</string>
  <key>StandardErrorPath</key>
  <string>__CACHE__/launchd-distill.err.log</string>
</dict>
</plist>
